# llm-proxy@.service
[Unit]
Description=LLM Proxy for Local LLM (%i mode)
After=network.target llama-server.service
Requires=llama-server.service

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_GROUP
WorkingDirectory=/path/to/llm_proxy
Environment="KILO_PROXY_LOG_DIR=/var/log/kilo_proxy"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/path/to/venv/bin/python /path/to/llm_proxy/__main__.py --mode %i
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
